# crontab
0 */1 * * * /usr/local/bin/run-refresher